# tgport.service
[Unit]
Description=tgport - Telegram bot wrapper for Claude CLI
After=network.target

[Service]
Type=simple
User=delta
WorkingDirectory=/home/delta/workspace/projects/tgport
EnvironmentFile=/home/delta/workspace/projects/tgport/.env
ExecStart=/home/delta/workspace/projects/tgport/.venv/bin/python -m tgport
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/delta/workspace
MemoryMax=512M

[Install]
WantedBy=multi-user.target
